# data/crontab
# /etc/crontab: system-wide crontab
# Unlike any other crontab you don't have to run the `crontab'
# command to install the new version when you edit this file
# and files in /etc/cron.d. These files also have username fields,
# that none of the other crontabs do.

# no user inline

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user  command

* * * * * php /var/www/cron.php > /var/www/cron.log 2>/dev/null

0 20 *  *  0 /usr/bin/certbot -auto renew 2> /dev/null > /var/www/certbot.log

1 1 1 * * echo . > /var/www/access.log 2>/dev/null

1 1 1 * * echo . > /var/www/error.log 2>/dev/null

1 1 1 * * echo . > /var/www/cron.log 2>/dev/null

# An empty line is required at the end of this file for a valid cron file.
